# systemd/greenhouse-webserver.service
[Unit]
Description=Greenhouse Webserver Service
After=network.target greenhouse-manager.service
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/greenhouse-pi
Environment="PATH=/home/pi/greenhouse-pi/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="FLASK_HOST=0.0.0.0"
Environment="FLASK_PORT=5000"
Environment="FLASK_DEBUG=False"

# Set authentication credentials (change these!)
Environment="GREENHOUSE_USERNAME=admin"
Environment="GREENHOUSE_PASSWORD=greenhouse"

# Run the webserver
ExecStart=/home/pi/greenhouse-pi/.venv/bin/python /home/pi/greenhouse-pi/src/webserver/app.py

# Restart configuration
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=greenhouse-webserver

# Resource limits
MemoryLimit=512M
CPUQuota=50%

[Install]
WantedBy=multi-user.target
